# rdm.service
[Unit]
Description=Launch InvenioRDM Web Interface

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/coda-test
ExecStart=/home/ubuntu/.local/bin/pipenv run uwsgi --logger syslog:inveniordm_ui uwsgi_ui.ini

[Install]
WantedBy=multi-user.target
